# Linux Admin/Day10-11/Fady Tasks/SystemD service/BackupService/backup.timer
[Unit]
Description=Daily Backup Timer

[Timer]
OnCalendar=*-*-* 20:09:00
Persistent=true
Unit=backup.service

[Install]
WantedBy=timers.target

#######WILL Be added to /etc/systemd/system#######